# files/var/spool/cron/crontabs/philbot
@reboot /home/philbot/supervise-philbot
