# files/var/spool/cron/crontabs/philbot
#@reboot /home/philbot/supervise-philbot
@reboot /usr/bin/screen -d -m bash -c '/home/philbot/launch-philbot; exec bash'
