# files/var/spool/cron/crontabs/searchbot
@reboot tmux new-session -d -s searchbotsession 'supybot /home/searchbot/searchbot.conf'
